# crontab
00 * * * * bundle exec rails search:twitter
