# crontab
00 */4 * * * bundle exec rails search:twitter
30 07 * * * bundle exec rails search:lists
*/15 * * * * bundle exec rails tokens:refresh
